# crontab/checkIllegal2.slurm
#!/bin/bash

# Let no files expand to a null string
shopt -s nullglob

# Go to the Appusage directory
cd /cluster/software/appusage/1.0

# Sort .illegal files and store them in appusage_log2
today=`date | cut -c5-10`
rm -f crontab/appusage_log2/*
for f in log/*.illegal
do
  # ignore compute nodes that are down, version 1
  lsout=`ls -l $f | grep "$today"`
  if [ "$lsout" != "" ]; then

    # ignore compute nodes that are down, version 2
    name=`basename $f`
    node=`echo $name | cut -d'.' -f1`
    nodeinfo=`sinfo -N | grep "$node " | tr '\n' ' '`
    if [ -n "$nodeinfo" ]; then

      # ignore compute nodes that are down, version 3
      down=`echo "$nodeinfo" | egrep 'down|drain|resv' | tr -d '[:space:]'`
      if [ -z "$down" ]; then
        sort $f > crontab/appusage_log2/$name
      fi
    fi
  fi
done

# Copy common lines from appusage_log1 and appusage_log2 to bin/msg
rm -f crontab/msg
touch crontab/msg
for f in crontab/appusage_log1/*
do
  name=`basename $f`
  if [ -e crontab/appusage_log2/$name ]; then
    comm -1 -2 $f crontab/appusage_log2/$name >> crontab/msg
  fi
done

# Send an e-mail with the content of bin/msg
msg=`cat crontab/msg 2>/dev/null`

if [ "$msg" != "" ]; then
  echo >> crontab/msg
  echo "This message was generated by Appusage." >> crontab/msg

  cat crontab/msg | mailx -s "Appusage: a process is hanging on your favorite computer" adm@somwhere.com
fi
